# minecraft-world-ramsync/src/crontab
*/5 * * * * /usr/bin/flock -n /dev/shm/worldsync.lock /App/save-world.sh
